# deploy/ollama.service.d/flash-attn.conf
# Ollama host-level binding for text-transform-service (cycle T13, 2026-07-14).
#
# WHY THIS EXISTS — opinion-gate classifies batches of up to ~100 candidate stories in a single
# call, which needs a large context window (TTS sizes num_ctx to the transform's budget: 14144 for
# opinion-gate). At that context size, qwen3.5:9b with an **f16** KV cache under flash attention
# closes the verdict array early under temp 0 — a SILENT tail-drop (e.g. 27/34 verdicts, no error),
# which the caller's fail-closed rule then turns into over-exclusion of the batch tail. Quantizing
# the KV cache to **q8_0** removes the instability: every batch (21/34/60) comes back complete.
#
# WHY BOTH VARS —
#   OLLAMA_KV_CACHE_TYPE=q8_0   the actual fix (stable large-context decoding + smaller KV cache).
#   OLLAMA_FLASH_ATTENTION=1    REQUIRED: llama.cpp refuses V-cache quantization without flash
#                               attention ("V cache quantization requires flash_attn" → segfault).
#                               The box already ran flash=auto→on, so this only pins it.
#
# This is a HOST-WIDE Ollama setting (affects every model/transform on the box), not a TTS config.
# It is safe for the other transforms: their contexts are small (KV cache is a few MiB either way)
# and q8_0 KV is an imperceptible quality change for these extraction tasks. See docs/models.md
# and the T13 CYCLE-LOG entry for the full evidence (VRAM: 100% GPU, ~4 GB headroom at 14336 ctx;
# 34-batch ~22 s; 3× stable).
#
# INSTALL (human, sudo — same boundary as deploy/redeploy.sh; see deploy/README.md §3a):
#   sudo mkdir -p /etc/systemd/system/ollama.service.d
#   sudo cp deploy/ollama.service.d/flash-attn.conf /etc/systemd/system/ollama.service.d/
#   sudo systemctl daemon-reload
#   sudo systemctl restart ollama
# Verify:  systemctl show ollama.service -p Environment   # shows both vars
#          # then a large opinion-gate batch returns a complete verdict set (id per input).

[Service]
Environment="OLLAMA_FLASH_ATTENTION=1"
Environment="OLLAMA_KV_CACHE_TYPE=q8_0"
